# 42/deployment/systemd/railway-monitor.service
[Unit]
Description=Railway Communication Monitor System
After=network.target mongodb.service redis.service

[Service]
Type=simple
User=railway
Group=railway
WorkingDirectory=/opt/railway-monitor
Environment=NODE_ENV=production
Environment=SERVER_PORT=3000
Environment=SIGNALING_PORT=8080
Environment=MONGODB_URI=mongodb://localhost:27017/railway_monitor
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=LOG_DIR=/var/log/railway-monitor
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
